# backend/deploy/tapten-backup.timer
[Unit]
Description=Run Tap Ten editorial backend backup daily

[Timer]
OnCalendar=*-*-* 03:25:00
Persistent=true
Unit=tapten-backup.service

[Install]
WantedBy=timers.target
